# infra/packer/ops/systemd/ops.target
[Unit]
Description=Ops Stack Entry Point
Wants=bootstrap.service
After=network-online.target

[Install]
WantedBy=multi-user.target
