# infra/packer/ops/systemd/platform-api.service
[Unit]
Description=Platform API
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
Restart=always
RestartSec=5
TimeoutStartSec=120

EnvironmentFile=/etc/platform/api.env

ExecStartPre=-/usr/bin/docker rm -f platform-api

ExecStartPre=/bin/sh -c '/usr/local/bin/aws ecr get-login-password --region us-east-1 \
  | /usr/bin/docker login --username AWS --password-stdin \
    123456789012.dkr.ecr.us-east-1.amazonaws.com'

ExecStartPre=/usr/bin/docker pull ${IMAGE_URI}

ExecStart=/usr/bin/docker run \
  --name platform-api \
  --rm \
  -p 127.0.0.1:${PORT}:3000 \
  -e NODE_ENV=${NODE_ENV} \
  ${IMAGE_URI}

ExecStop=/usr/bin/docker stop platform-api

[Install]
WantedBy=ops.target
